# associate/challenge_2/crontab
# CRON, what time is it?
0-59/1 * * * * /usr/bin/wattime.sh
*/10 * * * * python /usr/bin/dir_cleaner.py
